# fk-config-service-confd/etc/service/fk-config-service-confd/run
#!/bin/bash -e

exec 2>&1
exec /usr/share/fk-config-service-confd/bin/confd

